# target/linux/mediatek/mt7988/base-files/etc/init.d/fanspeed
#!/bin/sh /etc/rc.common

if [ ! -f /sys/class/pwm/pwmchip0/export ]; then
	exit 0
fi

